# backups/launchagents/20260405_030005_ai.openclaw.email-integration.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HARVEY_EMAIL_AUTH</key>
		<string>DSswQ3xnSWXgbkyK</string>
	</dict>
	<key>Label</key>
	<string>ai.openclaw.email-integration</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/fhjtech/.openclaw/workspace/email_integration/learning_integration.py</string>
		<string>--action</string>
		<string>process</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/havery_email.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/havery_email.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>9</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>
</dict>
</plist>
