# backups/launchagents/20260408_030001_com.hjtech.auto-learner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HARVEY_EMAIL_AUTH</key>
		<string>SEMefmThGnEKJiTz</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.hjtech.auto-learner</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/fhjtech/.openclaw/workspace/.scripts/auto_learner.py</string>
		<string>verify</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/fhjtech/.openclaw/logs/auto_learner.err</string>
	<key>StandardOutPath</key>
	<string>/Users/fhjtech/.openclaw/logs/auto_learner.log</string>
	<key>StartInterval</key>
	<integer>1800</integer>
</dict>
</plist>
